# filesystem/usr/lib/systemd/system/emuld.service
[Unit]
Description=emuld

[Service]
Type=simple
ExecStart=/usr/bin/emuld
