# filesystem/usr/lib/systemd/system/sensord.service
[Unit]
Description=sensord

[Service]
Type=simple
ExecStart=/usr/sbin/sensord
