# filesystem/etc/profile.d/emulator_ecore_workaround.sh
# W/A for incorrect rendering
export ECORE_EVAS_FORCE_SYNC_RENDER=1
