# filesystem/etc/profile.d/launcher.sh
# Supposedly used by menu-screen
export LAUNCHER_ENGINE=x11
