# filesystem/etc/profile.d/simulator-env.sh
export TSLIB_TSDEVICE=/dev/input/touchscreen0
export KEY_DEVICE=/dev/input/event0
export TSLIB_CONFFILE=/usr/etc/ts.conf
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/usr/lib:./lib:/opt/lib
export PATH=$PATH:/sbin:/usr/sbin:/opt/bin
export TSLIB_CALIBFILE=/usr/etc/pointercal
export ELM_FONT_PATH=/usr/share/SLP/fonts:/usr/share/SLP/licensed_fonts:/usr/share/fonts/truetype/ttf-bitstream-vera
export FONTNAME=HelveN_Asia
